# systemd/user/battery-notify.service
[Unit]
Description=Battery notification script

[Service]
Type=oneshot
ExecStart=%h/dotfiles/bin/battery_notify.sh
